# luci-app-backrest/root/etc/uci-defaults/luci-backrest
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@backrest[-1]
	add ucitrack backrest
	set ucitrack.@backrest[-1].init=backrest
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
